# swatch.target
[Unit]
Description=Start swatch sessions

[Install]
WantedBy=multi-user.target
